# supervisord.conf
[supervisord]
nodaemon=true

[program:docker]
command=/usr/sbin/nginx